# docs/ops/openrc/xp
#!/sbin/openrc-run

name="xp"
description="xp (Xray control plane)"

command="/usr/local/bin/xp"
command_args="run"
command_user="xp:xp"

# If you use an env file, keep it readable by root and the service user.
environment_file="/etc/xp/xp.env"

depend() {
  need net
  # Suggested if xray is managed on the same host:
  use xray
}

