# docs/ops/openrc/xp
#!/sbin/openrc-run

name="xp"
description="xp (Xray cluster manager)"

command="/usr/local/bin/xp"
command_args="run"
command_user="xp:xp"

depend() {
  need net
  # Suggested if xray is managed on the same host:
  use xray
}
